# snapgear-2.6-p42/vendors/gaisler/leon2mmu/romfs/etc/inittab
# Basic Init Table
#

# System initialization (runs when system boots).
::sysinit:/etc/init.d/rcS

# Get us a non-login interactive shell at a controlling tty
# with job control. ATM sh is a symlink to BASH.
# If any of busybox's shell are present, sh is a symlink to them
ttyS0::respawn:-/bin/sh
